# desktop/hugin-watch.path
[Unit]
Description=Hugin — watch for state changes and trigger notification
Documentation=https://github.com/sramkrishna/hugin-mcp-client

[Path]
PathModified=%h/.hugin/state/dashboard_state.json

[Install]
WantedBy=default.target
